# supervisord.conf
[supervisord]
nodaemon=true

[program:phddns]
command=/bin/bash -c "phddns start"

[program:status]
command=/bin/bash -c "cat /tmp/oraysl.status"
